# deploy/webtmux.service
[Unit]
Description=Webtmux - Mobile Web Terminal
After=network.target

[Service]
Type=simple
ExecStart=/data/develop/golang/webtmux/webtmux -listen :3400
Restart=on-failure
RestartSec=3
StandardOutput=journal
StandardError=journal

# Security hardening
NoNewPrivileges=true
ProtectSystem=strict
ProtectHome=read-only
ReadWritePaths=%h/.webtmux

# Environment (override password via drop-in or env file)
EnvironmentFile=-%h/.config/webtmux/env

[Install]
WantedBy=default.target
